# Bash/optimization-tools/files/50-spice-vdagent.rules
ACTION=="change",KERNEL=="card0", SUBSYSTEM=="drm", RUN+="/usr/local/bin/x-resize"
